# config/cron/crontab
0 */1 * * * ~/scripts-repo/scripts/check-updates.sh
